# pkg/rc.d/mdns_repeater
#!/bin/ksh

daemon="/usr/local/sbin/mdns-repeater"
daemon_flags="-p /var/run/mdns-repeater/mdns-repeater.pid"
daemon_user="_mdns-repeater"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
